# systemd/wg-frontend.service
[Unit]
Description=MikroTik WireGuard Manager Frontend
Documentation=https://github.com/mustafakiractr/wg-manager
After=network.target wg-backend.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/wg/frontend

# Production: Nginx dist/ klasörünü serve eder
# Bu servis sadece build update için kullanılır
# Normal durumda çalışmasına gerek yok (nginx yeterli)
ExecStart=/bin/true

# Otomatik yeniden başlatma
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Loglar
StandardOutput=append:/root/wg/frontend.log
StandardError=append:/root/wg/frontend_error.log

# Güvenlik ayarları
NoNewPrivileges=true
PrivateTmp=true

# Kaynak limitleri
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
